# files/etc/network/if-up.d/svxlink.restart
#!/bin/bash
echo $(date) >> /tmp/svxlink.start.ifup.txt
service svxlink restart

